# rootfs/lib/systemd/system/vm_config.target.wants/ssh-config.service
[Unit]
Description=Configure SSH
DefaultDependencies=no
#Requires=network-config.service
After=network-config.service

[Service]
ExecStart=/sbin/sshd -d
Type=oneshot
